# sysadmin/daily.timer
[Unit]
Description=Run Gamocosm daily processes

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
