# sysadmin/daily.service
[Unit]
Description=Gamocosm daily processes

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
EnvironmentFile=/root/gamocosm/gamocosm.env
ExecStart=/usr/bin/certbot renew
ExecStart=/usr/bin/systemctl restart nginx
ExecStart=/usr/bin/systemctl restart pod-gamocosm
ExecStart=/usr/bin/curl --silent --show-error --output /dev/null https://gamocosm.com
ExecStart=/usr/bin/bash /root/gamocosm/sysadmin/backup.sh
