# server_setup/zram.service
[Unit]
Description=Zram (source: Gamocosm)

[Service]
ExecStart=/usr/bin/zram-helper start

[Install]
WantedBy=multi-user.target
